# systemd/system/keepd-watcher.path
[Path]
PathChanged=/etc/keepd/policy.json

[Install]
WantedBy=multi-user.target
